# debian/easyroam-config.service
[Unit]
Description=Check your easyroam certificate and inform before expiration.

[Service]
Type=oneshot
ExecStart=/usr/bin/easyroam-config -w 30
